# deployment/init.d/simple-tftpd
#!/bin/sh
#
# simple-tftpd        Simple TFTP Daemon
#
# chkconfig: 2345 80 20
# description: Lightweight TFTP server (RFC 1350)
#
### BEGIN INIT INFO
# Provides:          simple-tftpd
# Required-Start:    $remote_fs $network $time
# Required-Stop:     $remote_fs $network $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Simple TFTP Daemon
# Description:       Lightweight TFTP server (RFC 1350)
### END INIT INFO

NAME=simple-tftpd
SERVICE="${NAME}.service"
DAEMON=/usr/bin/simple-tftpd
CONFIG=/etc/simple-tftpd/simple-tftpd.conf
PIDFILE=/run/simple-tftpd/simple-tftpd.pid
RUNAS=simple-tftpd

. /etc/rc.d/init.d/functions 2>/dev/null || . /lib/lsb/init-functions 2>/dev/null || true

if [ ! -x "$DAEMON" ]; then
    echo "$DAEMON not found or not executable" >&2
    exit 1
fi
if [ ! -f "$CONFIG" ]; then
    echo "Config file not found: $CONFIG" >&2
    echo "Copy an example from /etc/simple-tftpd/examples/ or run: sudo simple-tftpd --test-config" >&2
    exit 1
fi

systemd_unit_file() {
    [ -f "/usr/lib/systemd/system/${SERVICE}" ] \
        || [ -f "/lib/systemd/system/${SERVICE}" ] \
        || [ -f "/etc/systemd/system/${SERVICE}" ]
}

uses_systemd() {
    [ -x /usr/bin/systemctl ] && systemd_unit_file
}

# RHEL/CentOS/Rocky/Debian with systemd: delegate everything to systemctl
if uses_systemd; then
    case "$1" in
        start|stop|restart|status|condrestart|try-restart)
            exec /usr/bin/systemctl "$1" "$SERVICE"
            ;;
        force-reload|reload)
            exec /usr/bin/systemctl reload "$SERVICE"
            ;;
        *)
            echo "Usage: $0 {start|stop|restart|status|reload|force-reload}" >&2
            exit 1
            ;;
    esac
fi

# Non-systemd fallback (SysV only)
do_start() {
    mkdir -p "$(dirname "$PIDFILE")"

    if command -v start-stop-daemon >/dev/null 2>&1; then
        start-stop-daemon --start --quiet --oknodo --background --make-pidfile \
            --pidfile "$PIDFILE" --chuid "$RUNAS:$RUNAS" \
            --exec "$DAEMON" -- --foreground start --config "$CONFIG"
        return $?
    fi

    if command -v runuser >/dev/null 2>&1; then
        runuser -u "$RUNAS" -- "$DAEMON" --foreground start --config "$CONFIG" &
    else
        su -s /bin/sh "$RUNAS" -c "$DAEMON --foreground start --config $CONFIG" &
    fi
    echo $! >"$PIDFILE"
    return 0
}

do_stop() {
    if [ -f "$PIDFILE" ]; then
        kill "$(cat "$PIDFILE")" 2>/dev/null || true
        rm -f "$PIDFILE"
    fi
    pkill -x simple-tftpd 2>/dev/null || true
    return 0
}

case "$1" in
    start)
        do_start
        rc=$?
        if command -v success >/dev/null 2>&1; then
            [ "$rc" -eq 0 ] && success || failure
        elif [ "$rc" -ne 0 ]; then
            echo "Failed to start $NAME" >&2
            exit "$rc"
        else
            echo "$NAME started"
        fi
        ;;
    stop)
        do_stop
        rc=$?
        if command -v success >/dev/null 2>&1; then
            [ "$rc" -eq 0 ] && success || failure
        fi
        ;;
    restart|force-reload)
        do_stop
        do_start
        ;;
    status)
        if pgrep -x simple-tftpd >/dev/null 2>&1; then
            echo "$NAME is running"
            exit 0
        fi
        echo "$NAME is not running"
        exit 3
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|force-reload|status}" >&2
        exit 1
        ;;
esac

exit 0
