# deployment/systemd/simple-tftpd.service
[Unit]
Description=Simple TFTP Daemon
Documentation=file:/usr/share/doc/simple-tftpd/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simple-tftpd
Group=simple-tftpd
ExecStart=/usr/bin/simple-tftpd --foreground start --config /etc/simple-tftpd/simple-tftpd.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple-tftpd

# State and runtime directories (systemd creates under /var/lib, /run, /var/log)
StateDirectory=simple-tftpd
RuntimeDirectory=simple-tftpd
LogsDirectory=simple-tftpd

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/simple-tftpd /var/log/simple-tftpd /var/tftp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
