# deployment/systemd/{{PROJECT_NAME}}.service
[Unit]
Description=Simple TFTP Daemon
Documentation=file:/usr/share/doc/simple-tftpd/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simple-tftpd
Group=simple-tftpd
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/simple-tftpd --foreground start --config /etc/simple-tftpd/simple-tftpd.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple-tftpd

StateDirectory=simple-tftpd
RuntimeDirectory=simple-tftpd
LogsDirectory=simple-tftpd

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/simple-tftpd /var/log/simple-tftpd /var/tftp

[Install]
WantedBy=multi-user.target
