# packaging/linux/deb/postrm
#!/bin/bash
# Post-removal script for simple-tftpd (Debian)

set -e

PROJECT_NAME="simple-tftpd"
LOG_DIR="/var/log/${PROJECT_NAME}"

if [ "$1" = "purge" ]; then
    rm -rf "${LOG_DIR}"

    if command -v update-rc.d >/dev/null 2>&1; then
        update-rc.d "${PROJECT_NAME}" remove 2>/dev/null || true
    fi
fi

exit 0
